# MolClass/MolClassWeb/molclass/upstartMolClass.conf
description     "start fast MolClass structure search"

start on startup

task
exec /var/www/molclass/cron_MolClass.sh
